# ansible/remote/vms/storage/roles/setup_tape_archiver/templates/tape-archiver-tbot.service
# {{ ansible_managed }}
[Unit]
Description=tbot (Teleport Machine ID) for the tape-archiver data worker
Documentation=https://goteleport.com/docs/machine-workload-identity/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# tbot ships with the teleport-ent package already installed on this host.
# Runs as root: the TPM join needs /dev/tpmrm0 and the JWT destination is
# root-owned.
ExecStart=/usr/local/bin/tbot start -c /etc/tape-archiver/tbot.yaml --diag-addr 127.0.0.1:4444
# Readiness gate: tbot's /wait blocks until every configured service has
# reported its first status, and returns 200 only if all are healthy — the
# workload-identity-jwt service reports healthy strictly AFTER writing the
# JWT to disk. This makes the unit's start job complete only once the first
# JWT exists, so consumers get "ready", not just "process started", from a
# plain After=tape-archiver-tbot.service. On timeout/unhealthy the start
# fails and Restart=always retries the whole join.
# --retry-connrefused is load-bearing: ExecStartPost fires before tbot has
# bound the diag socket, and a refused connection otherwise fails curl
# INSTANTLY (exit 7) instead of waiting. --retry also treats an unhealthy
# 503 first report as transient.
ExecStartPost=/usr/bin/curl --silent --fail --output /dev/null --retry 30 --retry-connrefused --retry-delay 1 --retry-max-time 120 --max-time 120 http://127.0.0.1:4444/wait
# Headroom for the gate: the whole start job (join + first renewal + gate
# retries) must fit inside the start timeout or systemd kills it.
TimeoutStartSec=180
Restart=always
RestartSec=10
# JWT destination on tmpfs — the credential never touches persistent disk.
# Preserved across tbot restarts: the worker container bind-mounts this
# directory, and a bind mount pins the directory inode — if systemd removed
# and recreated it on every tbot stop, renewals would land in a new directory
# the running container can never see.
RuntimeDirectory=tape-archiver/jwt
RuntimeDirectoryMode=0700
RuntimeDirectoryPreserve=yes

[Install]
WantedBy=multi-user.target
